# scripts/server/systemd/codex-workspace-sync-update.timer
[Unit]
Description=Run Codex Workspace Sync GitHub updates hourly

[Timer]
OnCalendar=hourly
RandomizedDelaySec=10m
Persistent=true

[Install]
WantedBy=timers.target
