# scripts/server/systemd/codex-workspace-sync.service
[Unit]
Description=Codex Workspace Sync API
After=network.target

[Service]
Type=simple
WorkingDirectory=__ROOT__
Environment=PYTHONUNBUFFERED=1
Environment=CWS_APP_ROOT=__ROOT__
Environment=CWS_STATE_ROOT=__STATE_ROOT__
Environment=PORT=__PORT__
EnvironmentFile=-/etc/codex-workspace-sync/service.env
ExecStart=__PYTHON__ -m cws.server.bootstrap serve --host 0.0.0.0 --port __PORT__
Restart=always
RestartSec=5
User=__USER__
Group=__USER__

[Install]
WantedBy=multi-user.target
